# hermes.conf
; ================================================================
; Hermes Gateway Telegram — supervisord service
; ================================================================
; Set HERMES_CMD di Railway env vars dengan command lengkap:
;   contoh: HERMES_CMD=/usr/local/bin/hermes-gateway --config /etc/hermes/config.yaml
;
; Manage via SSH:
;   supervisorctl status hermes-gateway
;   supervisorctl start hermes-gateway
;   supervisorctl stop hermes-gateway
;   supervisorctl restart hermes-gateway
;   supervisorctl tail -f hermes-gateway
; ================================================================

[program:hermes-gateway]
command=%(ENV_HERMES_CMD)s
directory=/root
autostart=%(ENV_HERMES_AUTOSTART)s
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=10
stdout_logfile=/tmp/hermes.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=2
stderr_logfile=/tmp/hermes-error.log
stderr_logfile_maxbytes=2MB
stderr_logfile_backups=2
environment=HOME="/root",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
